# segmentinfo.timer
[Unit]
Description=Run DVB-C DOCSIS snooper segmentinfo

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
